# packages/cloud-builder/cloud_builder-0.2.56.tar.gz/cloud_builder-0.2.56/systemd/cb-collect.service
[Unit]
Description=CB Project Collect Service
After=syslog.target
Requires=network-online.target

[Service]
EnvironmentFile=-/etc/cloud_builder
Type=simple
Restart=on-failure
ExecStart=/usr/bin/cb-collect --project ${CB_PROJECT} --timeout ${CB_COLLECT_TIMEOUT} --update-interval ${CB_COLLECT_UPDATE_INTERVAL}
ExecStartPre=/bin/bash -c 'rm -f /var/lock/_projects_*'

[Install]
WantedBy=multi-user.target
